# Personal Projects/75-huawei.rules
KERNEL=="ttyUSB0", SUBSYSTEM=="tty", ACTION=="add", ATTRS{idVendor}=="12d1", ATTRS{idProduct}=="1506", RUN+="/usr/bin/net_up.sh"